# supervisord-capture-script-watcher.conf
[program:capture-script-watcher]
command=/capture-script-watcher.sh
priority=1000
